# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:queue-worker]
command=php /app/artisan queue:work --queue=default
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/queue-worker.log
numprocs=2
process_name=%(program_name)s_%(process_num)02d

[program:scheduler]
command=php /app/artisan schedule:work
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/storage/logs/scheduler.log
